# infrastructure/tradelayout-backtest-api.service
[Unit]
Description=TradeLayout Backtest API Server
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/projects/tradelayout-backtesting-engine
Environment=PATH=/home/ubuntu/venv/bin
Environment=CLICKHOUSE_HOST=localhost
Environment=CLICKHOUSE_PORT=8123
Environment=CLICKHOUSE_USER=default
Environment=CLICKHOUSE_PASSWORD=
Environment=CLICKHOUSE_DATABASE=tradelayout
Environment=CLICKHOUSE_TABLE=nse_ticks_indices
ExecStart=/home/ubuntu/venv/bin/python backtest_api_server.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
